# systemd/codex-telegram-bridge.service.example
[Unit]
Description=Codex Telegram Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/agent/codex-telegram-bridge
ExecStart=/usr/bin/node /home/agent/codex-telegram-bridge/src/bot.js
EnvironmentFile=/home/agent/codex-telegram-bridge/.env
Restart=on-failure
RestartSec=5
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/agent /tmp

[Install]
WantedBy=default.target
