# systemd/codex-telegram-bridge.service.example
[Unit]
Description=Codex Telegram Bridge
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/.codex-telegram-bridge
ExecStart=codex-telegram-bridge run
EnvironmentFile=%h/.codex-telegram-bridge/config.env
Restart=on-failure
RestartSec=5
NoNewPrivileges=false
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/home/agent /tmp

[Install]
WantedBy=default.target
